# radius.service
[Unit]
Description=Orbit frontend server

[Service]
WorkingDirectory=/var/orbit/radius/
ExecStart=/var/orbit/radius/radius.start
